# packages/sysutils/systemd/config/system.d/tmp-emulation.mount
[Unit]
Description=Emulation RAM disk
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=tmpfs
Where=/tmp/emulation
Options=size=256m
Type=tmpfs

[Install]
WantedBy=multi-user.target
