# packages/supervisedthinking/emulation/streaming/moonlight-qt/system.d/moonlight-qt.service
[Unit]
Description=Moonlight is an open source implementation of NVIDIAs GameStream
After=graphical.target
Requires=graphical.target

[Service]
ExecStart=/usr/bin/moonlight-qt.start
TimeoutStopSec=30
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
