# packages/supervisedthinking/service/rr-config-tool/source/system.d/service.rr-config-tool.service
[Unit]
Description=RR Config Tool
After=pulseaudio.service
Before=kodi.service
Wants=kodi.service pulseaudio.service

[Service]
ExecStart=/bin/sh /usr/share/kodi/addons/service.rr-config-tool/bin/rr-config-tool.start
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=kodi.service pulseaudio.service
