# packages/supervisedthinking/emulation/frontends/Emulationstation/emulationstation/system.d/emulationstation.service
[Unit]
Description=Emulation Station emulator front-end
After=graphical.target
Requires=graphical.target

[Service]
ExecStart=/usr/bin/emulationstation.start
TimeoutStopSec=30
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
